# crypto-bot.service
[Unit]
Description=Advanced Crypto Trading Bot
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=cryptobot
Group=cryptobot
WorkingDirectory=/opt/crypto_bot
Environment=PATH=/opt/crypto_bot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/crypto_bot/venv/bin/python3 bot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=30
StandardOutput=append:/opt/crypto_bot/logs/bot-stdout.log
StandardError=append:/opt/crypto_bot/logs/bot-stderr.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/crypto_bot/data /opt/crypto_bot/logs /opt/crypto_bot/models

# Memory limit for 4GB VPS
MemoryMax=2G
MemoryHigh=1536M

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
